# infrastructure/mesos-slave/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=0.0.0.0:9002

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mesos-slave]
command=mesos-slave --master=zk://%(ENV_ZOOKEEPER_HOSTS)s/mesos --containerizers=docker,mesos --executor_registration_timeout=5mins --hostname=%(ENV_MASTER_HOST)s --docker_stop_timeout=3secs
autorestart=true
user=root
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout] 
command = supervisor_stdout 
buffer_size = 100 
events = PROCESS_LOG 
result_handler = supervisor_stdout:event_handler
